# data/fantable.service.in
# This file is part of fantable
# More at https://github.com/everdrone/jetson-fan-table

[Unit]
Description=Fan Table Control
ConditionPathExists=@sysconfdir@/fantable/table
# Requires the GPU files to be set first
After=nvpmodel.service

[Service]
Type=simple
ExecStart=/usr/sbin/fantable
User=root

[Install]
WantedBy=multi-user.target
