# data/fantable.service.in
# this file is part of jetson-fantable
# more at https://github.com/everdrone/jetson-fantable

[Unit]
Description=Fan Table Control
ConditionPathExists=@sysconfdir@/fantable/table

[Service]
Type=simple
ExecStart=/usr/sbin/fantable
User=root

[Install]
WantedBy=multi-user.target
